# images/vnc/supervisor/conf.d/vncserver.conf
[program:vncserver]
command=/opt/vnc/scripts/vncserver.sh
autorestart=true
stopasgroup=true
killasgroup=true
priority=100
